# ospackage/svc/cryptctl2-client.service
[Unit]
Description=Disk encryption utility (cryptctl2) - key client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/cryptctl2 --action client-daemon
User=root
Group=root
WorkingDirectory=/
PrivateTmp=true
RestartSec=5
Restart=on-abort
MountFlags=shared

[Install]
WantedBy=multi-user.target
